# config/ihdas.service
[Unit]
Description=ihdas Simple URL Shortener
After=network.target postgresql.service

[Service]
Type=simple
User=ihdas
Group=ihdas
WorkingDirectory=/opt/ihdas
ExecStart=/opt/ihdas/ihdas
Restart=always
RestartSec=5

# Environment variables
Environment=DATABASE_URL=postgres://ihdas:secure-and-complicated-password@localhost/ihdas?sslmode=disable # Set your PostgreSQL connection string
Environment=PORT=8080 # App runs on port 8080
#
# # Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ihdas
#
[Install]
WantedBy=multi-user.target
